# root/etc/systemd/system/backuppc.service.d/nethserver.conf
[Service]
KillSignal=SIGTERM
Restart=always
RestartSec=1
