# root/etc/systemd/system/backuppc.service.d/nethserver.conf
[Service]
RuntimeDirectory=BackupPC
RuntimeDirectoryMode=0755
KillSignal=SIGTERM
Restart=always
RestartSec=1
